# bim-app.service
[Unit]
Description=BIM App Service
After=network.target

[Service]
WorkingDirectory=/root/bim_project
ExecStart=/bin/bash -c '[ -x /root/bim_project/.venv/bin/poetry ] && exec /root/bim_project/.venv/bin/poetry run bim-app-prod || exec /usr/local/bin/poetry run bim-app-prod'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target