# bim-app.service
[Unit]
Description=BIM App Service
After=network.target

[Service]
WorkingDirectory=/root/bim-project
ExecStart=/usr/bin/poetry run bim-app-prod
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target