# cloudflared.service
[Unit]
Description=Cloudflare Tunnel (simple mode)
After=network.target

[Service]
ExecStart=/usr/local/bin/cloudflared tunnel --url http://localhost:8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target